# deploy/systemd/ainews-api.service
[Unit]
Description=AI News & Trends Report — Uvicorn API
Documentation=https://github.com/NguyenSiTrung/AINewtrendsReport
After=network.target valkey-server.service
Wants=valkey-server.service

[Service]
Type=exec
User=ainews
Group=ainews
WorkingDirectory=/opt/ainews/app
EnvironmentFile=/etc/ainews/ainews.env
ExecStart=/opt/ainews/venv/bin/uvicorn ainews.api.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 2 \
    --log-level info
Restart=on-failure
RestartSec=5

# Security hardening
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/ainews /var/log/ainews

# Resource limits
LimitNOFILE=65536
MemoryMax=1G

[Install]
WantedBy=multi-user.target
